# systemd/smc-healthcheck.timer.template
[Unit]
Description=Run SMC Signal Engine health check every {{HEALTHCHECK_INTERVAL}}

[Timer]
OnBootSec=2min
OnUnitActiveSec={{HEALTHCHECK_INTERVAL}}
AccuracySec=30s
Persistent=true
Unit=smc-healthcheck.service

[Install]
WantedBy=timers.target
