# launchd/com.smc.signalengine.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.smc.signalengine</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/kanannagiev/Documents/New project/project/.venv/bin/python</string>
        <string>/Users/kanannagiev/Documents/New project/project/main.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/kanannagiev/Documents/New project/project</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONUNBUFFERED</key>
        <string>1</string>
    </dict>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/kanannagiev/Library/Logs/SMCSignalEngine/signal-engine.out.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/kanannagiev/Library/Logs/SMCSignalEngine/signal-engine.err.log</string>
</dict>
</plist>
